# bot.conf
[program:casting]
directory=/home/ubuntu/Casting
command=python3 app.py
autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/Casting/logfile_err.log
stdout_logfile=/home/ubuntu/Casting/logfile.log